# systemd/system/secure-time-sync.timer
[Unit]
Description=Sync time securely every day
After=network.target

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target
